# raspi-clock.conf
[program:raspi_clock]
command=python3 -u main.py
directory=/home/pi/raspi-clock
stdout_logfile=/home/pi/raspi-clock/output.txt
redirect_stderr=true