# raspi-clock.conf
[program:raspi_clock]
command=/usr/bin/python3 -u main.py
environment=HOME="/home/pi",USER="pi"
directory=/home/pi/raspi-clock/src
stdout_logfile=/home/pi/raspi-clock/output.txt
redirect_stderr=true